# lib/systemd/system/cron.service
[Unit]
Description=Regular background program processing daemon

[Service]
ExecStart=/usr/sbin/cron -f

[Install]
WantedBy=multi-user.target
